# fw-lid-fprint-daemon-1.0/laptop-lid-daemon.service
[Unit]
Description=Fingerprint toggle daemon (lid + externals)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/laptop-lid-daemon.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
